# cas-service.service
[Unit]
Description=CAS Microservice (SymPy + Maxima)
After=network.target

[Service]
Type=simple
User=sam
WorkingDirectory=/media/sam/1TB/cas-service
ExecStart=/home/sam/.local/bin/uv run python -m cas_service.main
Restart=on-failure
RestartSec=5
Environment=CAS_PORT=8769
Environment=CAS_MAXIMA_PATH=/usr/bin/maxima
Environment=CAS_LOG_LEVEL=INFO

[Install]
WantedBy=multi-user.target
