# cas-service.service
# Edit User, WorkingDirectory, and ExecStart to match your installation.
[Unit]
Description=CAS Microservice (Sage + SymPy + MATLAB + WolframAlpha)
After=network.target

[Service]
Type=simple
User=your-username
WorkingDirectory=/path/to/cas-service
ExecStart=/usr/local/bin/uv run python -m cas_service.main
Restart=on-failure
RestartSec=5
Environment=CAS_PORT=8769
Environment=CAS_LOG_LEVEL=INFO
# Environment=CAS_SAGE_PATH=sage
# Environment=CAS_SAGE_TIMEOUT=30
# Environment=CAS_MATLAB_PATH=matlab
# Environment=CAS_MATLAB_TIMEOUT=30
# Environment=CAS_WOLFRAMALPHA_APPID=  # get at developer.wolframalpha.com

[Install]
WantedBy=multi-user.target
